# novnc.conf
[program:novnc]
command=/usr/bin/novnc_server --listen 8080 --vnc localhost:5900
priority=90
stdout_logfile=/var/log/novnc.stdout.log
stderr_logfile=/var/log/novnc.stderr.log